# root/etc/s6-overlay/s6-rc.d/mrefd/run
#!/command/with-contenv sh

PIDFILE="/config/mrefd.pid"

# Create pid file for service uptime dashboard class
touch ${PIDFILE}

# start daemon
exec /usr/local/bin/mrefd ${MREFD_CONFIG_DIR}/mrefd.cfg
